# deploy/lux-mentor.service
[Unit]
# LETI Matematika Mentori — Telegram bot (systemd).
# DIQQAT: tavsiya etilgan yo'l — Docker Compose (deploy/DEPLOY_ORACLE_CLOUD.md).
# Bu service Docker'siz, to'g'ridan-to'g'ri venv bilan ishlatmoqchi bo'lganlar uchun
# muqobil. Bu holatda Redis alohida o'rnatilishi kerak (yoki REDIS_URL bo'sh = RAM fallback).
Description=LETI Matematika Mentori (Telegram bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/leti-ai-mentor
ExecStart=/home/ubuntu/leti-ai-mentor/.venv/bin/python run_mentor.py

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=leti-mentor

[Install]
WantedBy=multi-user.target
